# robot-harness/deploy/robot-harness.service
[Unit]
Description=Onchain Rover Rust robot harness
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=@ROVER_ENV_FILE@
WorkingDirectory=@ROVER_HARNESS_DIR@
ExecStartPre=/bin/sh -lc "pgrep -f '[p]ython.*(app.py|api:app|read_serial|capture_images|voice)|[u]vicorn.*api:app' | xargs -r kill"
ExecStartPre=/bin/sh -lc "sleep 0.5; pgrep -f '[p]ython.*(app.py|api:app|read_serial|capture_images|voice)|[u]vicorn.*api:app' | xargs -r kill -9"
ExecStartPre=/bin/sh -lc "pgrep -f '[r]over-harness($| )' | xargs -r kill"
ExecStart=@ROVER_BIN@
ExecStop=/bin/sh -lc "curl -fsS --max-time 2 -X POST http://127.0.0.1:8000/motors/stop >/dev/null 2>&1 || true"
Restart=always
RestartSec=2
KillSignal=SIGINT
KillMode=mixed
TimeoutStartSec=20
TimeoutStopSec=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
